# hedge-fund-core/src/main/resources/upstart/hedge-historian.conf
start on started networking

script
    java -cp /opt/hedge-fund/${env.CI_COMMIT_ID}/lib/hedge-fund-*.jar com.onenow.main.HistorianMain
end script